# demo-1/services/turqette/llama-qwen27-4090.service
[Unit]
Description=llama.cpp Qwen3.5-27B Q4_K_M on RTX 4090
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=CUDA_VISIBLE_DEVICES=0
ExecStart=llama-server -m /path/to/models/Qwen3.5-27B-GGUF/Qwen3.5-27B-Q4_K_M.gguf --jinja --chat-template-file /path/to/models/qwen3.5_chat_template.jinja -ngl 99 -c 131072 -np 1 -fa auto -ctk q4_0 -ctv q4_0 -sm none -mg 0 --fit off --mmproj /home/turq/models/Qwen3.5-27B-GGUF/mmproj-F16.gguf --host 0.0.0.0 --port 8080
Restart=always
RestartSec=5
TimeoutStartSec=0

[Install]
WantedBy=default.target
