# demo-4/services/pg1/llama-demo4-122b-q6-split.service
[Unit]
Description=Demo4 llama.cpp Qwen3.5-122B UD-Q6_K_XL split on RTX PRO 6000 + RTX 5090 (pg1)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/turq
Environment=CUDA_VISIBLE_DEVICES=0,1
ExecStart=/home/turq/src/llama.cpp/build-dualarch/bin/llama-server -m /home/turq/models/Qwen3.5-122B-A10B-GGUF/UD-Q6_K_XL/Qwen3.5-122B-A10B-UD-Q6_K_XL-00001-of-00004.gguf --jinja --chat-template-file /home/turq/models/qwen3.5_chat_template.jinja -ngl 99 -c 131072 -np 1 -fa auto -ctk q4_0 -ctv q4_0 -sm layer -ts 3,1 --fit off --reasoning-format none --reasoning-budget 0 --host 0.0.0.0 --port 18084
Restart=always
RestartSec=5
TimeoutStartSec=0

[Install]
WantedBy=default.target
